# tests/tools/systemd/dbus.socket
[Unit]
Description=D-Bus System Message Bus Socket

[Socket]
ListenStream=/shared/dbus/system_bus_socket
